# packaging/app-mail-report.cron
20 4 * * *   root   /usr/sbin/pflogsumm /var/log/maillo*[g1-4] > /var/clearos/mail_report/data-month.out
25 4 * * *   root   /usr/sbin/pflogsumm -d yesterday /var/log/maillo*[g1] > /var/clearos/mail_report/data-yesterday.out
27 * * * *   root   /usr/sbin/pflogsumm -d today /var/log/maillo*[g1] > /var/clearos/mail_report/data-today.out
